# packages/cool-booy-tmux/etc/profile.d/cool-booy-tmux.sh
#!/bin/bash
if ! command -v tmux &>/dev/null; then
    echo "Tmux is not installed. Installing tmux..."
    apt-get install tmux
fi
